# systemd_units/ark-core.service
[Unit]
Description=ARK Core Service
Documentation=https://github.com/arklinux
After=network.target
Wants=ark-watchdog.service

[Service]
Type=simple
User=ark
Group=ark
WorkingDirectory=/opt/ark
ExecStart=/usr/bin/python /opt/ark/agents/core.py
Restart=on-failure
RestartSec=5s

# Security hardening
NoNewPrivileges=yes
PrivateTmp=yes
PrivateDevices=yes
ProtectSystem=strict
ProtectHome=yes
ReadOnlyPaths=/opt/ark/models
ReadWritePaths=/opt/ark/logs /opt/ark/memory /opt/ark/bus /opt/ark/ingest /opt/ark/run
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictNamespaces=yes
RestrictRealtime=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM
IPAddressAllow=localhost
IPAddressDeny=any

# Resource limits
MemoryMax=2G
CPUQuota=80%
TasksMax=256

[Install]
WantedBy=multi-user.target
